# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
loglevel=info

; ── Synthia Node.js MCP bus ───────────────────────────────────────────────────
[program:synthia-node]
command=node --max-old-space-size=1024 -r ./server/control-preload.js server/lite.js
directory=/app/synthia-server
autostart=true
autorestart=true
startretries=5
startsecs=5
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=
    NODE_ENV="%(ENV_NODE_ENV)s",
    PORT="%(ENV_SYNTHIA_NODE_PORT)s",
    NODE_MODE="%(ENV_NODE_MODE)s",
    SUPABASE_PRIMARY_URL="%(ENV_SUPABASE_PRIMARY_URL)s",
    SUPABASE_PRIMARY_KEY="%(ENV_SUPABASE_PRIMARY_KEY)s",
    SUPABASE_SECONDARY_URL="%(ENV_SUPABASE_SECONDARY_URL)s",
    SUPABASE_SECONDARY_KEY="%(ENV_SUPABASE_SECONDARY_KEY)s",
    TERMINAL_TOKEN="%(ENV_TERMINAL_TOKEN)s",
    CORS_ORIGIN="%(ENV_CORS_ORIGIN)s",
    ANTHROPIC_API_KEY="%(ENV_ANTHROPIC_API_KEY)s",
    OPENAI_API_KEY="%(ENV_OPENAI_API_KEY)s",
    GITHUB_TOKEN="%(ENV_GITHUB_TOKEN)s",
    HF_TOKEN="%(ENV_HF_TOKEN)s",
    SYNTHIA_URL="%(ENV_SYNTHIA_URL)s",
    TRIDENT_ONNX_ENABLED="%(ENV_TRIDENT_ONNX_ENABLED)s",
    TRIDENT_MODEL_URL="%(ENV_TRIDENT_MODEL_URL)s",
    KIMI_API_KEY="%(ENV_KIMI_API_KEY)s"

; ── Synthia Python FastAPI server ────────────────────────────────────────────
[program:synthia-python]
command=python server.py
directory=/app/synthia-server
autostart=true
autorestart=true
startretries=3
startsecs=10
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=
    PORT="%(ENV_SYNTHIA_PYTHON_PORT)s",
    ANTHROPIC_API_KEY="%(ENV_ANTHROPIC_API_KEY)s",
    OPENAI_API_KEY="%(ENV_OPENAI_API_KEY)s",
    SUPABASE_URL="%(ENV_SUPABASE_PRIMARY_URL)s",
    SUPABASE_KEY="%(ENV_SUPABASE_PRIMARY_KEY)s",
    KIMI_API_KEY="%(ENV_KIMI_API_KEY)s"

; ── Studio Express server (serves built React frontend) ──────────────────────
[program:studio]
command=node dist/index.js
directory=/app/studio
autostart=true
autorestart=true
startretries=5
startsecs=5
priority=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=
    NODE_ENV="%(ENV_NODE_ENV)s",
    PORT="%(ENV_PORT)s",
    SUPABASE_URL="%(ENV_SUPABASE_PRIMARY_URL)s",
    SUPABASE_KEY="%(ENV_SUPABASE_PRIMARY_KEY)s",
    GITHUB_TOKEN="%(ENV_GITHUB_TOKEN)s",
    GOOGLE_ACCESS_TOKEN="%(ENV_GOOGLE_ACCESS_TOKEN)s",
    GOOGLE_CLIENT_ID="%(ENV_GOOGLE_CLIENT_ID)s",
    GOOGLE_CLIENT_SECRET="%(ENV_GOOGLE_CLIENT_SECRET)s",
    ANTHROPIC_API_KEY="%(ENV_ANTHROPIC_API_KEY)s",
    OPENAI_API_KEY="%(ENV_OPENAI_API_KEY)s",
    KIMI_API_KEY="%(ENV_KIMI_API_KEY)s",
    SYNTHIA_URL="%(ENV_SYNTHIA_URL)s",
    PYTHON_API_URL="%(ENV_PYTHON_API_URL)s",
    PYTHON_API_AUTOSTART="%(ENV_PYTHON_API_AUTOSTART)s",
    LINUX_CONTAINER_ENABLED="%(ENV_LINUX_CONTAINER_ENABLED)s",
    LINUX_CONTAINER_WORKDIR="%(ENV_LINUX_CONTAINER_WORKDIR)s",
    LINUX_CONTAINER_SHELL="%(ENV_LINUX_CONTAINER_SHELL)s",
    MESH_DATA_DIR="%(ENV_MESH_DATA_DIR)s",
    TERMINAL_TOKEN="%(ENV_TERMINAL_TOKEN)s"
